# install/upstart.conf
description "ppagent"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid ppagent
setgid ppagent
chdir {chdir}

exec {exec_path} -c /etc/ppagent/config.json run
